# etc/demo-api.service
[Unit]
Description=Demo API Server

[Service]
WorkingDirectory=/path/to/api/directory
ExecStart=/path/to/cargo run --release

[Install]
WantedBy=multi-user.target
